# misc/systemd/startup.service
[Unit]
Description=startup

[Service]
Type=oneshot
ExecStart=/root/arch/scripts/startup.sh

[Install]
WantedBy=multi-user.target
