# crontabs/igor
*/5 * * * *    /home/igor/arch/scripts/dns.sh
*/5 * * * *    /home/igor/arch/scripts/updates.sh
