# crontabs/root
0 2 * * 0      /home/igor/arch/scripts/fstrim.sh
