# crontabs/root
0 */6 * * *    /home/igor/arch/scripts/fstrim.sh
* * * * *      /home/igor/arch/scripts/cifswatch.sh
