# cron/crontab
*/5 * * * *    /home/igor/arch/scripts/dns.sh
* * * * *      /home/igor/arch/scripts/dpms.sh
*/5 * * * *    /home/igor/arch/scripts/updates.sh
* * * * *      /home/igor/arch/scripts/df.sh