# linode/crontab
0 1 * * *     /home/igor/certbot/renew.sh
