# system/crontab
@reboot   /home/igor/arch/scripts/updates.sh
0 * * * * /home/igor/arch/scripts/updates.sh
* * * * * /home/igor/arch/scripts/dns.sh