# system/crontab
@reboot   root /home/igor/arch/scripts/updates.sh
* * * * * root /home/igor/arch/scripts/updates.sh
* * * * * igor /home/igor/arch/scripts/dns.sh