# linode/openrc/nodeping
#!/sbin/openrc-run
 
user="nodey"
group="nodey"
command="/usr/bin/node"
command_args="/var/www/ping/ping.js"
command_user="${user}:${group}"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
 
depend() {
	use net
}
