# linode/openrc/nodecontact
#!/sbin/openrc-run

user="nodey"
group="nodey"
command="/usr/bin/node"
command_args="/var/www/contactform/src/main.js"
command_user="${user}:${group}"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"

export NODE_ENV="production"

depend() {
        use net
}
